# mysql.conf
[program:mysql]
user = root
datadir = /app/mysql
port = 3306
log-bin = /app/mysql/mysql-bin

command=mysqld_safe

autorestart=true
exitcodes=0

stopsignal=TERM
stopasgroup=true
stopwaitsecs=10

stdout_events_enabled=true
stderr_events_enabled=true
